# Systemd/loraaprs.timer
[Unit]
Description=LoraAPRS-GW Service

[Timer]
OnBootSec=2min

[Install]
WantedBy=multi-user.target



